# internal/ignition/templates/systemd-journal-gatewayd.socket
[Unit]
Description = Fake systemd-journal-gatewayd.socket

[Socket]
ListenStream = 19531
Accept = yes

[Install]
WantedBy = sockets.target
